# xos/synchronizers/vnfm/supervisor/vtn-observer.conf
[program:vtn-observer]
command=python /opt/xos/observers/vbng/vtn-observer.py -C /opt/xos/observers/vbng/vtn_observer_config
