# xos/observers/hpc/supervisor/hpc-watcher.conf
[program:hpc-watcher]
command=python /opt/xos/observers/hpc/hpc-watcher.py
