# xos/observers/hpc/supervisor/hpc-observer.conf
[program:hpc-observer]
command=python /opt/xos/observers/hpc/hpc-observer.py -C /opt/xos/observers/hpc/hpc_observer_config
