# observer.conf
[program:observer]
command=python /opt/xos/xos-observer.py
