# xos/onboard/vsg/synchronizer/supervisor/vcpe-observer.conf
[program:vcpe-observer]
command=python /opt/xos/synchronizers/vsg/vcpe-synchronizer.py -C /opt/xos/synchronizers/vsg/vcpe_synchronizer_config
