# xos/openstack_observer/templates/container.service.j2
[Unit]
Description={{ container_name }}
After=docker.service

[Service]
ExecStart=/bin/bash -c "/usr/local/sbin/start-container-{{ container_name }}.sh"

[Install]
WantedBy=multi-user.target
