# xos/observers/vcpe/supervisor/vcpe-observer.conf
[program:vcpe-observer]
command=python /opt/xos/observers/vcpe/vcpe-observer.py -C /opt/xos/observers/vcpe/vcpe_observer_config
