# xos/synchronizers/vbng/supervisor/vbng-observer.conf
[program:vbng-observer]
command=python /opt/xos/observers/vbng/vbng-observer.py -C /opt/xos/observers/vbng/vbng_observer_config
