# xos/synchronizers/monitoring_channel/supervisor/monitoring_channel_observer.conf
[program:monitoring_channel_observer]
command=python /opt/xos/observers/monitoring_channel/monitoring_channel_observer.py -C /opt/xos/observers/monitoring_channel/monitoring_channel_observer_config
