# xos/synchronizers/vcpe/supervisor/vcpe-observer.conf
[program:vcpe-observer]
command=python /opt/xos/synchronizers/vcpe/vcpe-synchronizer.py -C /opt/xos/synchronizers/vcpe/vcpe_synchronizer_config
