# xos/observers/onos/supervisor/onos-observer.conf
[program:onos-observer]
command=python /opt/xos/observers/onos/onos-observer.py -C /opt/xos/observers/onos/onos_observer_config
